# packages/Debian/deenesse/etc/systemd/system/deenesse.service
[Unit]
Description=Update DNS on Cloudflare

[Service]
ExecStart=/usr/bin/deenesse /etc/deenesse/%I

[Install]
WantedBy=network.target
